# ½ű¾/ÅäÖÃÎļþ½âÎö/99-usb.rules
SUBSYSTEMS=="usb",ACTION=="add",KERNEL=="sd[a-z]*",ATTRS{serial}=="4C531001550521103504",SYMLINK+="mingda-usb%n",RUN+="/sbin/service httpd restart"
